# etc/supervisor/conf.d/syzygy.conf
[program:syzygy]
command=bash -c "ulimit -n 4096; exec python3 /srv/syzygy-tables.info/server.py"
user=www-data
autostart=true
directory=/srv/syzygy-tables.info/
