# etc/supervisor/conf.d/syzygy.conf
[program:syzygy]
command=bash -c "ulimit -n 4096; source venv/bin/activate; python server.py"
user=www-data
autostart=true
directory=/home/niklas/syzygy-tables.info/
