# etc/systemd/system/syzygy.service
[Unit]
Description=Syzygy tablebases
After=network.target
Requires=network.target

[Service]
User=www-data
WorkingDirectory=/home/niklas/syzygy-tables.info/
ExecStart=/bin/bash -c "ulimit -n 4096; source venv/bin/activate; exec python server.py"
Restart=always
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
